# scraper/systemd/instagram-scraper.timer.example
[Unit]
Description=Run IG Dashboard Instagram scraper daily

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
Unit=instagram-scraper.service

[Install]
WantedBy=timers.target
