# scraper/systemd/reddit-scraper.timer.example
[Unit]
Description=Run IG Dashboard Reddit scraper daily

[Timer]
OnCalendar=*-*-* 06:15:00
Persistent=true
Unit=reddit-scraper.service

[Install]
WantedBy=timers.target
